# supervisor-app.conf
[program:app-uwsgi]
command = /home/docker/code/uwsgi-2.0.11.2 /home/docker/code/uwsgi_conf.xml

[program:nginx-app]
command = /usr/sbin/nginx
